# media/jellyfin/rclone-homelaptop-watchdog.service
[Unit]
Description=Watchdog for rclone homelaptop mount

[Service]
Type=oneshot
ExecStart=/srv/mystackmint/scripts/rclone-watchdog.sh
