# media/jellyfin/rclone-homelaptop.service
[Unit]
Description=rclone mount of home laptop G: drive
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=notify
ExecStartPre=-/bin/fusermount -uz /mnt/homelaptop
ExecStartPre=/bin/mkdir -p /mnt/homelaptop
ExecStart=/usr/bin/rclone mount homelaptop:/G: /mnt/homelaptop \
  --config /root/.config/rclone/rclone.conf \
  --allow-other --read-only \
  --vfs-cache-mode full --vfs-cache-max-size 5G \
  --timeout 10s --contimeout 10s --low-level-retries 3 --retries 1 \
  --sftp-idle-timeout 1m --dir-cache-time 30s --poll-interval 0
ExecStop=/bin/fusermount -uz /mnt/homelaptop
Restart=always
RestartSec=15

[Install]
WantedBy=multi-user.target
